# writeboost
#!/bin/sh

### BEGIN INIT INFO
# Provides:          writeboost
# Required-Start:    $local_fs
# Required-Stop:     $local_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: create and remove a writeboost device
### END INIT INFO

# --------------------------------------------------

# (mandatory)

BACKING=/dev/hvg2/back18g # VM
# BACKING=/dev/mapper/pdc_fgaahfadb # server

CACHE=/dev/hvg1/cache2g # VM
# CACHE=/dev/sdb1 # server

PLOG=/dev/hvg1/cache3m # VM
# PLOG=/dev/sdb2 # server

# --------------------------------------------------

# (optional)

VOLNAME=writeboost-vol # create /dev/mapper/writeboost-vol
MNTDIR=/exports/writeboost-vol

# --------------------------------------------------

MODPATH=/lib/modules/`uname -r`/kernel/drivers/md

do_start() {
    # load the required modules
    modprobe libcrc32c
    insmod ${MODPATH}/dm-writeboost.ko

    # create the device
    dmsetup create ${VOLNAME} --table "0 `blockdev --getsize ${BACKING}` writeboost 1 ${BACKING} ${CACHE} ${PLOG} 4 segment_size_order 10 nr_rambuf_pool 8 6 enable_migration_modulator 1 sync_interval 0 update_record_interval 0"

    mount /dev/mapper/${VOLNAME} ${MNTDIR}
}

do_stop() {
    umount -l ${MNTDIR}

    # better to record the last migrated id to shorten the reboot time
    dmsetup message ${VOLNAME} 0 update_record_interval 1

    # explicitly suspend/resume for stable shutdown
    dmsetup suspend ${VOLNAME}
    dmsetup resume ${VOLNAME}

    dmsetup remove ${VOLNAME}
}

case "$1" in
    start)
        do_start
        ;;
    stop)
        do_stop
        ;;
    status)
        dmsetup status ${VOLNAME}
        ;;
    restart|reload|force-reload)
        ;;
    *)
        echo "Usage: writeboost {start|stop|status}" >&2
        ;;
esac
